# Icap/icap.cron
# Run ICAP test every 1 minute
*/1 * * * * root /home/admin/Moneta/icap-test.sh
